# deploy/host/shorts-batch.timer
[Unit]
Description=nightly batch trigger

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
